# tests/test_docker/nc-client/config/states/500-nc-client/supervisor.conf
[program:nc-client]
command=wait-for-tcp -h nc-server -p 4242 -- nc-client
autostart=true
priority=500
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=6
